# config/bigsdb/cron/crontab/bigsdb_crontab
# The following jobs are added (i.e. appended) to /etc/crontab
* * * * *	bigsdb	xvfb-run -a /usr/local/bin/bigsjobs
05 * * * *	bigsdb	/usr/local/bin/sync_user_dbase_users.pl --user-database $BIGSDB_USER_DATABASE
00 18 * * 0 	postgres	psql -c "DELETE FROM guid WHERE last_accessed < NOW() - INTERVAL '6 months'" bigsdb_prefs
00 5 * * *	postgres psql -c "DELETE FROM jobs where (stop_time IS NOT NULL AND stop_time < now()-interval '7 days') OR (status LIKE 'rejected%' AND submit_time < now()-interval '7 days')" bigsdb_jobs > /dev/null 2>&1
